# status-report/status-report.service
[Unit]
Description=Status reporting service

[Service]
Type=oneshot
ExecStart=/usr/share/status-report/status-report.sh
User=root
