# status-report/status-report.timer
[Unit]
Description=Status reporting timer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
