# hosts-blocklists-dnsmasq/fetch-dns-blacklist.timer
[Unit]
Description=Update dnsmasq tracking blacklist daily

[Timer]
OnCalendar=*-*-* 5:00:00
Persistent=true

[Install]
WantedBy=timers.target
