# gsm-alerts/gsm-alerts.service
[Unit]
Description=Gsm alerting service

[Service]
Type=oneshot
ExecStart=/usr/share/gsm-alerts/gsm-alerts.sh
User=root
